# supervisor/sam.conf
[program:sam_test]
#process_name=%(program_name)s_%(process_num)02d
#numprocs=2
command=python3 ./test.py
directory=/home/hebron/sam
autostart=false
autorestart=true
user=hebron
